# oled-display.service
[Unit]
Description=OLED Display for Serveo URL
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/selmi/WifiVisualizer
ExecStart=/usr/bin/python3 /home/selmi/WifiVisualizer/serveo_oled.py
StandardOutput=append:/var/log/oled-display.log
StandardError=append:/var/log/oled-display-error.log
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
